# scripts/com.brain.transcript-snapshotter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brain.transcript-snapshotter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>BRAIN_SCRIPTS="${BRAIN_SCRIPTS:-$HOME/brain-template/scripts}"; "$BRAIN_SCRIPTS/snapshot-transcripts.sh" "${BRAIN_ROOT:-$HOME/brain}"</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/tmp/brain-snapshotter.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/brain-snapshotter.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
